# etc/init.d/iptables
#!/bin/sh
### BEGIN INIT INFO
# Provides:       iptables
# Required-Start:
# Required-Stop:
# Default-Start:
# Default-Stop:
# Description:    Stub — iptables is not used in this container
### END INIT INFO

# LSB exit codes for "status":
#   0 = running
#   3 = program is not running
# SUM's firewall detection calls "service iptables status" and treats
# anything other than exit 3 as "firewall is running".

case "$1" in
    status)
        echo "iptables: not running (stub)"
        exit 3
        ;;
    *)
        exit 0
        ;;
esac
